# example/docker/spva_ldap/ldap-supervisor.conf
[program:ldap]
command=/usr/sbin/slapd -h "ldap:///" -d 1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ldap.out.log
stderr_logfile=/var/log/supervisor/ldap.err.log
