# example/docker/spva_ldap/sssd-supervisor.conf
[program:sssd]
command=/usr/sbin/sssd -i
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sssd.out.log
stderr_logfile=/var/log/supervisor/sssd.err.log
