# example/docker/spva_krb/krb5-kdc-supervisor.conf
[program:krb5-kdc]
command=/usr/sbin/krb5kdc -n
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/krb5kdc.err.log
stdout_logfile=/var/log/supervisor/krb5kdc.out.log
