# example/docker/spva_jwt/jwt-issuer-supervisor.conf
[program:jwt-issuer]
command=python3 /opt/epics/app.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/jwt-issuer.out.log
stderr_logfile=/var/log/supervisor/jwt-issuer.err.log
