# example/docker/spva_krb5/krb5-admin-server-supervisor.conf
[program:krb5-admin-server]
command=/usr/sbin/kadmind -nofork
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/kadmind.err.log
stdout_logfile=/var/log/supervisor/kadmind.out.log
